# mydbus.service
[Unit]
Description=Inhibitor test
After=sysinit.target

[Service]
ExecStart=/mnt/local11/dbus/mydbus

[Install]
WantedBy=default.target
